# hetzner-slack-bot.service
[Unit]
Description=Hetzner Slack Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hetzner-slack-bot
ExecStart=/root/.nvm/versions/node/v24.14.0/bin/node app.js
EnvironmentFile=/opt/hetzner-slack-bot/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
